# sabnzbd.conf
[program:sabnzbd]
environment = LANG="en_GB.UTF-8",LC_ALL="en_GB.UTF-8"
user = nobody
command = /usr/sbin/python2 /opt/sabnzbd/SABnzbd.py --config-file /config --server 0.0.0.0:8080 --https 8090